# etc/supervisor/conf.d/marathon.conf
[program:marathon]
command = /usr/bin/marathon
directory = /var/lib/mesos
redirect_stderr = true
