# puppet/modules/consular/templates/init.consular.conf.erb
description "consular"

start on stopped rc RUNLEVEL=[2345]
respawn

exec /usr/bin/consular<%
  @consular_args.each do |arg| %> <%= arg %><% end %>
